# varnish-ssh/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:varnish]
command=varnishd -f /etc/varnish/custom.vcl -s malloc,100M -a 0.0.0.0:80
